# roles/komodo_periphery/templates/periphery.service.j2
# Managed by Ansible — do not edit manually.
# Source: roles/komodo_periphery/templates/periphery.service.j2

[Unit]
Description=Komodo Periphery Agent v{{ komodo_version }}
Documentation=https://komo.do/docs
After=network-online.target docker.service
Wants=network-online.target
# Soft dependency on Docker — periphery still starts if Docker is absent.

[Service]
Type=simple
User={{ periphery_user }}
Group={{ periphery_group }}
ExecStart={{ periphery_binary_path }} --config-path {{ periphery_config_path }}
Restart=on-failure
RestartSec=5s
# Give periphery time to gracefully shut down before killing it.
TimeoutStopSec=30s

# Harden the service slightly.
PrivateTmp=true
NoNewPrivileges=true
# The role creates the service user and adds it to the 'docker' group
# automatically (periphery_manage_user: true). Disable if managing the
# user externally, or set periphery_user: root to skip user management.

[Install]
WantedBy=multi-user.target
